# some_code_examples_from_SVN/IPFIX-to-ClickHouse/etc/cron.d/rugion-jq2ch
*/2 * * * * root /usr/share/jq2ch/json2ch.sh
5 * * * * root /usr/share/jq2ch/ch_aggregate.sh







