# some_code_examples_from_SVN/IPFIX-to-ClickHouse/etc/init/docker-jq2ch.conf
description "jquery to clickhouse"
author "Aleksey Stolyarov"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
normal exit 0 1 TERM HUP

pre-start script
    /usr/bin/docker rm -f jquery2ClickHouse 2>/dev/null || true

    if grep -qs '/mnt/volumes/IPFIX2ClickHouse' /proc/mounts; then
       echo "/mnt/volumes/IPFIX2ClickHouse already mounted."
    else
        modprobe zram num_devices=1 
        echo 2048M > /sys/block/zram0/disksize 
        mkfs.ext2 /dev/zram0 
        mount /dev/zram0 /mnt/volumes/IPFIX2ClickHouse
    fi

end script

script
    docker run --rm  --name jquery2ClickHouse  \
    -e TZ="Asia/Yekaterinburg"  \
    -v /etc/ipfixcol/:/etc/ipfixcol/  \
    -v /mnt/volumes/IPFIX2ClickHouse/:/tmp/ipfixcol/flow/ \
    --net=host  cesnet:ipfixcol -c /etc/ipfixcol/startup-json.xml
end script


post-stop script
    /usr/bin/docker stop jquery2ClickHouse
end script
